# Dockerfiles/appserver/supervisord/appserver-php5-fpm.conf
[program:appserver-php5-fpm]
command=/opt/appserver/sbin/php-fpm --fpm-config /opt/appserver/etc/php-fpm.conf
autorestart=true
priority=9050
