# Dockerfiles/appserver/supervisord/appserver.conf
[program:appserver]
command=/opt/appserver/bin/php -dappserver.php_sapi=appserver /opt/appserver/server.php
autorestart=true
user=www-data
priority=9000
