# Dockerfiles/appserver/supervisord/appserver-watcher.conf
[program:appserver-watcher]
command=/opt/appserver/bin/php -dappserver.php_sapi=appserver /opt/appserver/server.php -w
autorestart=true
priority=9010
