# Dockerfiles/icinga2/supervisord/icinga_apache2_reload.conf
[program:icinga_apache2_reload]
command=/bin/bash -c "sleep 5 && /etc/init.d/apache2 reload"
priority=9999
